# .config/systemd/user/emacs.service
[Service]
Type=simple
Restart=always
ExecStart=emacs --fg-daemon