# .config/systemd/user/emacs.service
[Unit]
Description=Emacs server

[Service]
Type=simple
Restart=always
ExecStart=emacs --fg-daemon